# .test/quagga/supervisord.conf
[supervisord]
nodaemon=true

[program:zebra]
command=/usr/sbin/zebra -f /etc/quagga/zebra.conf

[program:bgp]
command=/usr/sbin/bgpd -f /etc/quagga/bgpd.conf

